# bin/opt/guybrush/guybrush.service
[Unit]
Description=Guybrush - An asistant for Telegram
After=network.target

[Service]
Type=simple
ExecStart=/bin/bash /opt/guybrush/run.sh

[Install]
WantedBy=multi-user.target