# bin/DEBIAN/postinst
service=guybrush

addgroup ${service}
useradd --system --shell /bin/false -g ${service} ${service}
systemctl enable /opt/${service}/${service}.service
systemctl start ${service}
ln -s /etc/nginx/sites-available/${service} /etc/nginx/sites-enabled/${service}