# bin/DEBIAN/prerm
service=guybrush

systemctl --quiet is-active ${service} 
if [ "$?" = "0" ]; then
    systemctl stop ${service}
    systemctl disable ${service}
fi
